# vagrant/salt/root/webserver/upstart/apache.conf
description "apache2"
start on vagrant-mounted

script
    exec service apache2 start
end script